# tools/linux/41-usb-sas.rules
# USB-SAS
SUBSYSTEM=="usb", ACTION=="add", ATTRS{idVendor}=="04b4", ATTRS{idProduct}=="4001", RUN+="/bin/sh /bin/usb_sas $parent"

